# provisioning/minio/files/minio.service
[Unit]
Description=Minio server

[Service]
User=vagrant
Group=vagrant
ExecStart=/home/vagrant/minio server /vagrant/storage
